# docker/supervisord-hq.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord_hq.log
pidfile=/tmp/supervisord_hq.pid
childlogdir=/tmp

[program:streamlit_hq]
directory=/home/user/app
# [수정] 절대경로를 지우고 시스템이 PATH에서 찾도록 수정했습니다.
command=streamlit run app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true --server.maxUploadSize=50 --server.enableCORS=false --server.enableXsrfProtection=false --logger.level=info
autostart=true
autorestart=true
stdout_logfile=/tmp/streamlit_hq.out.log
stderr_logfile=/tmp/streamlit_hq.err.log

[program:uvicorn_api]
directory=/home/user/app
environment=PYTHONPATH="/home/user/app"
# [수정] 절대경로를 지우고 시스템이 PATH에서 찾도록 수정했습니다.
command=uvicorn modules.hq_analysis_api:app --host 0.0.0.0 --port 18080 --log-level info
autostart=true
autorestart=true
stdout_logfile=/tmp/uvicorn_hq.out.log
stderr_logfile=/tmp/uvicorn_hq.err.log

[program:nginx_front]
# [수정] -g "daemon off;" 옵션을 추가하여 Nginx가 배경으로 숨지 못하게 고정합니다.
command=/usr/sbin/nginx -g "daemon off;" -c /home/user/app/nginx-hq.conf
autostart=true
autorestart=true
priority=20
startsecs=2
stdout_logfile=/tmp/nginx_hq_wrap.out.log
stderr_logfile=/tmp/nginx_hq_wrap.err.log
